# runtime/scheduled-tasks/RG-card-lint.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.3" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <Description>Weekly Sunday 23:30 鈥?run /card-lint --mode=hygiene autonomously (audit 02 Cards/ graph: orphans, stale, broken links, frontmatter drift)</Description>

    <URI>\RG-card-lint</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-3323412737-987631274-1346458851-500</UserId>

      <LogonType>InteractiveToken</LogonType>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>

    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <RestartOnFailure>

      <Count>1</Count>

      <Interval>PT15M</Interval>

    </RestartOnFailure>

    <StartWhenAvailable>true</StartWhenAvailable>

    <IdleSettings>

      <Duration>PT10M</Duration>

      <WaitTimeout>PT1H</WaitTimeout>

      <StopOnIdleEnd>false</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>

  </Settings>

  <Triggers>

    <CalendarTrigger>

      <StartBoundary>2026-05-23T23:30:00+08:00</StartBoundary>

      <ScheduleByWeek>

        <WeeksInterval>1</WeeksInterval>

        <DaysOfWeek>

          <Sunday />

        </DaysOfWeek>

      </ScheduleByWeek>

    </CalendarTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>{{USER_HOME}}\.claude\skills\card-lint\run_wrapper.cmd</Command>

    </Exec>

  </Actions>

</Task>

