# runtime/scheduled-tasks/RG-reminder-feishu-gate.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <URI>\RG-reminder-feishu-gate</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-3323412737-987631274-1346458851-500</UserId>

      <LogonType>InteractiveToken</LogonType>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <IdleSettings>

      <Duration>PT10M</Duration>

      <WaitTimeout>PT1H</WaitTimeout>

      <StopOnIdleEnd>true</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

  </Settings>

  <Triggers>

    <TimeTrigger>

      <StartBoundary>2026-07-06T09:00:00+08:00</StartBoundary>

    </TimeTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>C:\Program Files\Git\usr\bin\bash.exe</Command>

      <Arguments>-lc {{VAULT_ROOT}}/.claude/_state/reminder-gate-promo-feishu.sh</Arguments>

    </Exec>

  </Actions>

</Task>

