# runtime/scheduled-tasks/RG-wechat-bridge.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-06-05T21:10:00</Date>
    <Author>{{HOSTNAME}}\Administrator</Author>
    <URI>\RG-wechat-bridge</URI>
    <Description>Native WeChat 4.0 reader bridge (WeFlow replacement). Serves WeFlow's local API on :5031 from direct SQLCipher-4 decrypt. Disabled by default — enable at cutover (微信 open, WeFlow.exe stopped).</Description>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>{{HOSTNAME}}\Administrator</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <Enabled>false</Enabled>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Triggers>
    <LogonTrigger>
      <StartBoundary>2026-06-05T21:10:00</StartBoundary>
    </LogonTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>powershell</Command>
      <Arguments>-NoProfile -WindowStyle Hidden -ExecutionPolicy Bypass -File "{{USER_HOME}}\Developer\operator-harness\windows\wechat\run-bridge.ps1"</Arguments>
    </Exec>
  </Actions>
</Task>
